# scripts/kaede/stork-kaede.service
# Stork 2.0 kæde-kurér (gov-5, plan V21 step 10) — systemd USER-unit.
# Install (manuelt, del af aktiverings-tjeklisten — IKKE auto-aktiveret):
#   mkdir -p ~/.config/systemd/user
#   cp scripts/kaede/stork-kaede.service ~/.config/systemd/user/
#   loginctl enable-linger $USER   # services overlever session-luk (preflight tjekker)
#   systemctl --user daemon-reload && systemctl --user enable --now stork-kaede
# Stop = manuelt flow (krav 7): systemctl --user stop stork-kaede
[Unit]
Description=Stork 2.0 kaede-kurer (gov-5 dirigent)
After=network-online.target
# Start-loft (rette-til punkt 11a, rodårsag #147): selv hvis en restart-vej
# overses, starter unit'en max 5 gange pr. 10 min — derefter failed-state.
StartLimitIntervalSec=600
StartLimitBurst=5

[Service]
Type=simple
WorkingDirectory=%h/stork-2.0
# PATH eksplicit: user-units arver ikke shell-PATH (claude/gh). Node + codex
# pinnes IKKE til en nvm-versions-sti (rette-til punkt 5/5c) — de afledes af
# .nvmrc via node-env.sh (sources i ExecStart + preflight), så nvm-opgradering
# ikke dræber kæden midt i natlig drift.
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStartPre=/bin/bash scripts/kaede/preflight.sh
ExecStart=/bin/bash -c 'source scripts/kaede/node-env.sh && exec node scripts/kaede/dirigent.mjs'
Restart=on-failure
RestartSec=30
# Verdikt-exits genopliver ALDRIG (punkt 11a): 2 = KAEDE-STOP/PAUSE-verdikt,
# 64 = guard-nægtelse (lås/baseline), 78 = persistent stop-fil (preflight +
# dirigent nægter). Stop-tilstanden bæres af scripts/kaede/.kaede-stop —
# fjernelse er Mathias' genåbnings-handling (aktiverings-tjekliste).
RestartPreventExitStatus=2 64 78

[Install]
WantedBy=default.target
